# postgres/supervisord.conf
[supervisord]
nodaemon=false

[program:pgagent]
command=pgagent -f -l=2 hostaddr=127.0.0.1 dbname=postgres user=postgres
autorestart=true
user=postgres